# docker/php/schedule/crontab
* * * * * /usr/local/bin/php /var/www/app/artisan schedule:run >> /dev/stdout 2>&1
